# backend/butiks-api.service
[Unit]
Description=Butiks API Backend Server
Documentation=https://github.com/aleynatila/butiks
After=network.target mongodb.service

[Service]
Type=simple
User=root
WorkingDirectory=/butiks/backend
Environment=NODE_ENV=production
Environment=PATH=/root/.nvm/versions/node/v20.19.6/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/root/.nvm/versions/node/v20.19.6/bin/node src/server.js
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=butiks-api

[Install]
WantedBy=multi-user.target
